# micro/cloudpods/mpanel-cloudpods.service
# Example systemd unit: CloudPods Provisioner Worker

[Unit]
Description=MigraCloud CloudPods Provisioner Microservice
After=network-online.target postgresql.service redis.service
Wants=network-online.target

[Service]
Type=simple
User=mhadmin
WorkingDirectory=/opt/micro/cloudpods
EnvironmentFile=/opt/micro/cloudpods/.env
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudpods-worker

[Install]
WantedBy=multi-user.target
